# files/sentry-worker.service
[Unit]
Description = Sentry Worker
After = network.target

[Service]
ExecStart = sentry run worker

[Install]
WantedBy = multi-user.target
