# files/sentry-web.service
[Unit]
Description = Sentry Web
After = network.target

[Service]
ExecStart = sentry run web

[Install]
WantedBy = multi-user.target
